# dist/_lib/systemd/system/heartbeat.service
[Unit]
Description=Heartbeat
After=network.target

[Service]
ExecStart=/bin/startheart
Restart=always

[Install]
Alias=heartbeat.service
