# dist/systemd/heartbeat.service
[Unit]
Description=Heartbeat
After=network.target

[Service]
ExecStart=heartbeat
Restart=always

[Install]
Alias=heartbeat.service
